# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:toolbox]
command=/usr/local/bin/toolbox --config /workspace/db/tools.yaml --address 127.0.0.1 --port 5000
autostart=true
autorestart=true
startsecs=2
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:app]
command=/bin/sh -c 'until python -c "import socket; s=socket.socket(); s.settimeout(1); s.connect((\"127.0.0.1\", 5000)); s.close()" >/dev/null 2>&1; do echo "Waiting for toolbox on 127.0.0.1:5000..."; sleep 1; done; uv run uvicorn app.api.main:app --host "${FASTAPI_HOST:-0.0.0.0}" --port "${PORT:-8080}"'
autostart=true
autorestart=true
startsecs=2
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
